# 0x0F-load_balancer/1-install_load_balancer
#!/usr/bin/env bash
# Script to install and configure HAProxy Load balancer

# Update the system
sudo apt update
sudo apt upgrade -y

# Install HAProxy
sudo apt install -y haproxy

# Configure HAProxy
cat << EOF | sudo tee /etc/haproxy/haproxy.cfg
global
    maxconn 4096

defaults
    mode http
    timeout connect 5s
    timeout client 50s
    timeout server 50s

frontend http_front
    bind *:80
    default_backend http_back

backend http_back
    balance roundrobin
    server web-01 3.86.7.18:80 check
    server web-02 34.207.190.90:80 check
EOF

# Replace 'web-01-ip' and 'web-02-ip' with the actual IP addresses of web-01 and web-02

# Create an init script for HAProxy
cat << 'EOF' | sudo tee /etc/init.d/haproxy
#!/bin/sh
### BEGIN INIT INFO
# Provides:          haproxy
# Required-Start:    $local_fs $network
# Required-Stop:     $local_fs $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: HAProxy init script
# Description:       HAProxy is a TCP/HTTP reverse proxy which is particularly suited for high availability environments.
### END INIT INFO

PATH=/sbin:/bin:/usr/sbin:/usr/bin
DAEMON=/usr/sbin/haproxy
CONFIG=/etc/haproxy/haproxy.cfg
PIDFILE=/var/run/haproxy.pid

# Additional options that are passed to the Daemon.
EXTRAOPTS=""

test -x $DAEMON || exit 0

[ -f /etc/default/haproxy ] && . /etc/default/haproxy

. /lib/lsb/init-functions

haproxy_start() {
    $DAEMON -f "$CONFIG" -D -p $PIDFILE $EXTRAOPTS
    sleep 2
    pidofproc -p $PIDFILE $DAEMON > /dev/null
}

haproxy_stop() {
    pidofproc -p $PIDFILE $DAEMON > /dev/null && haproxy_reload
    pidofproc -p $PIDFILE $DAEMON > /dev/null && killproc -p $PIDFILE $DAEMON -USR1
    pidofproc -p $PIDFILE $DAEMON > /dev/null || killproc -p $PIDFILE $DAEMON
}

haproxy_reload() {
    pidofproc -p $PIDFILE $DAEMON > /dev/null && killproc -p $PIDFILE $DAEMON -USR2
    pidofproc -p $PIDFILE $DAEMON > /dev/null || return 2
    return 0
}

haproxy_status() {
    status_of_proc -p $PIDFILE $DAEMON haproxy
}

haproxy_configtest() {
    $DAEMON -c -q -V -f "$CONFIG"
    if [ $? -eq 0 ]; then
        echo "Configuration file is valid"
    else
        echo "Error in configuration file, check with 'haproxy -c -f $CONFIG'"
    fi
}

haproxy_force_reload() {
    pidofproc -p $PIDFILE $DAEMON > /dev/null && killproc -p $PIDFILE $DAEMON -USR2
    return 0
}

case "$1" in
    start)
        haproxy_start
        ;;
    stop)
        haproxy_stop
        ;;
    reload|force-reload)
        haproxy_reload
        ;;
    restart)
        haproxy_stop
        sleep 1
        haproxy_start
        ;;
    status)
        haproxy_status
        ;;
    configtest)
        haproxy_configtest
        ;;
    *)
        echo "Usage: $0 {start|stop|reload|force-reload|restart|status|configtest}"
        exit 1
        ;;
esac
exit 0
EOF

# Make the script executable
sudo chmod +x /etc/init.d/haproxy

# Start and enable the HAProxy service
sudo service haproxy start
sudo systemctl enable haproxy

# Print success message
echo "HAProxy has been installed and configured successfully."
